# systemd/menu_bot.service
[Unit]
Description=InnerTrade Menu Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@APP_USER@
Group=@APP_GROUP@
WorkingDirectory=@APP_ROOT@
Environment=PYTHONUNBUFFERED=1
ExecStart=@APP_ROOT@/.venv/bin/python -u menu_bot.py
Restart=always
RestartSec=2
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
